# gislab-unit/static/cron.daily/fstrim
#!/bin/sh
LOG=/var/log/fstrim.log
echo "*** $(date -R) ***" >> $LOG
fstrim -v / >> $LOG
