# system/server/107-service-worker/conf/serf/serf.conf
description "Serf agent"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [06]

kill signal INT    # Use SIGINT instead of SIGTERM so serf can depart the cluster.
respawn            # Restart the process if it dies and GOAL was not 'stopping'.
kill timeout 90    # Allow 90 seconds for serf to die before sending SIGKILL.

exec /usr/local/bin/serf agent \
	-bind $(hostname  -I | awk -F" " '{print $NF}') \
	-event-handler "member-leave,member-failed=/usr/local/bin/serf-member-left.sh" \
	-tag role=worker \
	-syslog

pre-stop script
	/usr/local/bin/serf leave
	sleep 5
end script

post-stop script
	sleep 5
end script
