# system/roles/client-web/templates/init/gislab-web.conf.j2
description "GIS.lab Web"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec {{ root_dir }}/www/gunicorn.sh
