# gislab-unit/templates/cron.daily/fstrim
#!/bin/bash

LOG=/var/log/fstrim.log
echo "*** $(date -R) ***" >> $LOG

fstrim --verbose / >> $LOG
fstrim --verbose /storage >> $LOG
