# system/roles/service-worker/templates/serf/serf-join.conf.j2
description "Join a GIS.lab network"

start on runlevel [2345]
stop on runlevel [!2345]

task
respawn

script
sleep 5
exec /usr/local/bin/serf join {{ GISLAB_NETWORK_SERVER_IP }}
end script
